# roles/backup-client/files/backup-db.service
[Unit]
Description=Backup etc catalog 

[Service]
Type=simple
ExecStart=/opt/backup.sh 
